# initscripts/tssrvl.conf
description "HP TouchPad landscape touchscreen service"
start on (stopped tssrvp or started udev)

respawn
exec /usr/sbin/tssrvl
