# initscripts/tssrvp.conf
description "HP TouchPad portrait touchscreen service"
start on stopped tssrvl
stop on starting lightdm

respawn
exec /usr/sbin/tssrvp

